# deploy/systemd/secure-platform-api.service
[Unit]
Description=Secure License Platform API
After=network.target

[Service]
WorkingDirectory=/opt/secure-license-platform
EnvironmentFile=-/etc/secure-platform/api.env
ExecStart=/usr/bin/dotnet /opt/secure-license-platform/publish/api/Platform.Api.dll
Restart=always
RestartSec=5
Environment=ASPNETCORE_ENVIRONMENT=Production

[Install]
WantedBy=multi-user.target
