# deploy/systemd/secure-platform-admin.service
[Unit]
Description=Secure License Platform Admin
After=network.target

[Service]
WorkingDirectory=/opt/secure-license-platform
EnvironmentFile=-/etc/secure-platform/admin.env
ExecStart=/usr/bin/dotnet /opt/secure-license-platform/publish/admin/Platform.Admin.dll
Restart=always
RestartSec=5
Environment=ASPNETCORE_ENVIRONMENT=Production

[Install]
WantedBy=multi-user.target
